# scripts/com.user.tiddlywiki.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.user.tiddlywiki.sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>rsync -av --delete "$HOME/dev/tiddlywiki-starter/output/" "$HOME/Library/Mobile Documents/com~apple~CloudDocs/tiddlywiki/"</string>
  </array>
  <key>StartInterval</key><integer>600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>$HOME/Library/Logs/tw-sync.log</string>
  <key>StandardErrorPath</key><string>$HOME/Library/Logs/tw-sync.err</string>
</dict>
</plist>